# mesh-kindle-chat/rootfs/etc/services.d/meshchat/run
#!/usr/bin/with-contenv bashio
exec httpd -f -p 8765 -h /var/www
